# crontab
# http://d.hatena.ne.jp/anatoo/20120424/1335198888
* * * * * /usr/local/bin/growlnotify -m "姿勢正しく"
